# xu4fanctl-1/DEBIAN/postinst
#!/bin/sh
echo "Enabling fanctl..."
systemctl enable fanctl.service
systemctl start fanctl.service
