# xu4fanctl-1/DEBIAN/prerm
#!/bin/sh
echo "Disabling fanctl..."
systemctl stop fanctl.service
systemctl disable fanctl.service
