# xu4fanctl-1/etc/systemd/system/fanctl.service
[Unit]
Description=Manage odroid xu4 fan

[Service]
Type=simple
ExecStart=/sbin/fanctl.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
